# 云服务器上的数据/upload/cloud/imagedownload.service
[Unit]
Description=python images download program
[Service]
Type=forking
ExecStart=/usr/bin/sh /www/wwwroot/image.haoxx.top/public/startdownload.sh
[Install]
WantedBy=multi-user.target
